# BackupDns15-04-17/var/lib/dpkg/info/apt-show-versions.conffiles
/etc/cron.daily/apt-show-versions
/etc/apt/apt.conf.d/20apt-show-versions
/etc/bash_completion.d/apt-show-versions
